# conf/supervisor.conf
[program:server-run]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hw21/current/app/app.php server
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
startsecs=1
redirect_stderr=true
stopwaitsecs=3600
stderr_logfile=/var/log/supervisor/test.err.log
stdout_logfile=/var/log/supervisor/test.out.log