# supervisor/consumer.conf
[program:order_bank_statement]
command=php /data/application.local/bin/console rabbitmq:consumer -m 100 order_bank_statement --env=dev -vv
process_name=order_bank_statement_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/application.local/var/log/supervisor.order_bank_statement.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/application.local/var/log/supervisor.order_bank_statement.error.log
stderr_capture_maxbytes=1MB
